# recipes-connectivity/bluez5/files/serial-btattach@.timer
[Unit]
Description=Timer for Bluetooth attach service on %I

[Timer]
OnStartupSec=6

[Install]
WantedBy=timers.target
